# etc/cron.d/lucas_project.template
# Cron schedule mirroring orchestrator jobs
0 * * * * lucas_user /usr/bin/python /opt/lucas/main.py
